# rootfs/etc/update-motd.d/00-header
#!/bin/sh
printf "Welcome to SSH Server\n"
